# ShatterMe/Agent/shatterme.service
[Unit]
Description=ShatterMe

[Service]
Type=oneshot
ExecStart=/usr/bin/start

[Install]
WantedBy=multi-user.target
